# launchd/com.xuying.daily-posts.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xuying.daily-posts.update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/xuying/playground/daily-posts/scripts/update_post_xml.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/xuying/playground/daily-posts</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/xuying/playground/daily-posts/logs/launchd.update_post_xml.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/xuying/playground/daily-posts/logs/launchd.update_post_xml.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
